# BangumiBot.service

[Unit]

Description=Bangumi Bot

[Service]

ExecStart=node /home/harvey/Bangumi_Bot/index.js
